# avcMiner.service
[Unit]
Description=Webchain
After=network.target
[Service]
Type=forking
GuessMainPID=no
ExecStart= webchain -c ~./webMiner/config.json -l ~./webMiner/avcminer.log -B
Restart=always
User=root
[Install]
WantedBy=multi-user.target
